# anonimvaqti-bot.service
[Unit]
Description=AnonimVaqti Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AnonimVaqtiBot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/AnonimVaqtiBot/venv/bin/python3 bot.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/anonimvaqti-bot.log
StandardError=append:/var/log/anonimvaqti-bot.error.log

[Install]
WantedBy=multi-user.target 